# deploy/supervisor/clipforge-youtube.conf
[program:clipforge-youtube]
directory=/www/wwwroot/clipper/clipforge
command=/www/wwwroot/clipper/clipforge/scripts/recreate-compose-up.sh --watch-chrome
autostart=true
autorestart=true
startsecs=20
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/clipforge-youtube-supervisor.log
stderr_logfile=/tmp/clipforge-youtube-supervisor.err.log
environment=HOME="/home/fcn88",DISPLAY=":0",XAUTHORITY="/home/fcn88/.Xauthority"
